# systemd/vault.service.d/agent.conf.example

[Service]
ExecStart=/usr/local/bin/vault agent -config=/etc/vault.d/

# vim: ft=systemd
